# legacy/.Trash-1010/files/pkg-plist.3
Error: Orphaned: etc/X11/xinit/xinitrc.d/10-qt5-check-opengl2.sh
Error: Orphaned: usr/lib/.build-id/00/6ed027db327f827de94ac5b120b68178cc02e8
Error: Orphaned: usr/lib/.build-id/07/df77b55158e74042aa38d14abdc620a3cba341
Error: Orphaned: usr/lib/.build-id/13/3d098ed4fe87817fdc3dc8cc4cef72f7a7790c
Error: Orphaned: usr/lib/.build-id/16/a8c0d36307eb16a2dbfd40f2f41b2e12dfb494
Error: Orphaned: usr/lib/.build-id/21/e7b66f7156a6d1101cf364370363e46184f8c6
Error: Orphaned: usr/lib/.build-id/2e/cbde321a71ab77658644e914a371640228fc76
Error: Orphaned: usr/lib/.build-id/35/a84cb8188a0dac6192579442c554a82dc7e864
Error: Orphaned: usr/lib/.build-id/35/d1ef2b7da0379043b49aafbf9b522aa8b2ca6d
Error: Orphaned: usr/lib/.build-id/40/4a5cc9450360f458f7967ae8f158d6238a5cad
Error: Orphaned: usr/lib/.build-id/49/e8a0823ebc43533c9b7c423f68b8e5458c3d39
Error: Orphaned: usr/lib/.build-id/4f/783091a1e3ee588783e21c80bb5d918894d70c
Error: Orphaned: usr/lib/.build-id/57/3343e0324fa85a841b2b12d6995992c9249da1
Error: Orphaned: usr/lib/.build-id/5a/9b3788ac240549ad87318b61c14c35e779172b
Error: Orphaned: usr/lib/.build-id/61/3b48ffd78b0477ac5f963a6536c8eba14ce8cf
Error: Orphaned: usr/lib/.build-id/66/8d665b4d0e420e3655fa3a02f95cbc32a06cd3
Error: Orphaned: usr/lib/.build-id/66/ccdcc7082a750276114858bde65b916aeb7792
Error: Orphaned: usr/lib/.build-id/6d/714719f828d2e829bf57a3a6a60e22bcb27e13
Error: Orphaned: usr/lib/.build-id/8a/1043ca842a3c0c0401f22777b67e251b0b3220
Error: Orphaned: usr/lib/.build-id/8b/0f30d684764fd0af00075dc5ae875de3ebfef1
Error: Orphaned: usr/lib/.build-id/9e/d335188d9a62d8c640e10bbd916b767d000c42
Error: Orphaned: usr/lib/.build-id/ae/7a787ef0bf1d986e500a4d69e26c68e432e66a
Error: Orphaned: usr/lib/.build-id/b6/665cd00626592b21840892b4edaee74ccc17b2
Error: Orphaned: usr/lib/.build-id/bc/0a81c156d8e5bfec49a6591ed0598c5f98a127
Error: Orphaned: usr/lib/.build-id/bc/fe8d1e45ef31443b5e9ca1d42bec97f9e4b523
Error: Orphaned: usr/lib/.build-id/be/2524eca048425afb37564437710c79ee4ebe29
Error: Orphaned: usr/lib/.build-id/be/3cfcc2ef59fb311211338c8fdd55a5e69ad6df
Error: Orphaned: usr/lib/.build-id/c3/892e308eee025399842115d1a266559479fcce
Error: Orphaned: usr/lib/.build-id/c4/c836b3fbe37ddc4413b75412ab48fc83438080
Error: Orphaned: usr/lib/.build-id/c8/2483792a4cf278866c11cdbd809e834b2bd945
Error: Orphaned: usr/lib/.build-id/c8/d2c64af076be38a399b2586130bfec9e02110b
Error: Orphaned: usr/lib/.build-id/c9/2cdc6b516d08aa2979bd0344eeaff32ebd3a86
Error: Orphaned: usr/lib/.build-id/d8/dc3ad9a00e024a2a27f31d710629c3b7ff624f
Error: Orphaned: usr/lib/.build-id/e9/7db79fd017087f5379007f3e2776ae77227af9
Error: Orphaned: usr/lib/.build-id/f2/80eeba326b08af4d6a55b21999a531a58251ed
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QComposePlatformInputContextPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSEmulatorIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSKmsEglDeviceIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSKmsGbmIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEglFSX11IntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevKeyboardPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevMousePlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevTabletPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QEvdevTouchScreenPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QGifPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QGtk3ThemePlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QICOPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QIbusPlatformInputContextPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QJpegPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QLibInputPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QLinuxFbIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QMinimalEglIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QMinimalIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QOffscreenIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QTuioTouchPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QVncIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbEglIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbGlxIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QXcbIntegrationPlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5Gui/Qt5Gui_QXdgDesktopPortalThemePlugin.cmake
Error: Orphaned: usr/lib64/cmake/Qt5PrintSupport/Qt5PrintSupport_QCupsPrinterSupportPlugin.cmake
Error: Orphaned: usr/lib64/libQt5EglFSDeviceIntegration.so.5
Error: Orphaned: usr/lib64/libQt5EglFSDeviceIntegration.so.5.15
Error: Orphaned: usr/lib64/libQt5EglFSDeviceIntegration.so.5.15.9
Error: Orphaned: usr/lib64/libQt5EglFsKmsSupport.so.5
Error: Orphaned: usr/lib64/libQt5EglFsKmsSupport.so.5.15
Error: Orphaned: usr/lib64/libQt5EglFsKmsSupport.so.5.15.9
Error: Orphaned: usr/lib64/libQt5Gui.so.5
Error: Orphaned: usr/lib64/libQt5Gui.so.5.15
Error: Orphaned: usr/lib64/libQt5Gui.so.5.15.9
Error: Orphaned: usr/lib64/libQt5OpenGL.so.5
Error: Orphaned: usr/lib64/libQt5OpenGL.so.5.15
Error: Orphaned: usr/lib64/libQt5OpenGL.so.5.15.9
Error: Orphaned: usr/lib64/libQt5PrintSupport.so.5
Error: Orphaned: usr/lib64/libQt5PrintSupport.so.5.15
Error: Orphaned: usr/lib64/libQt5PrintSupport.so.5.15.9
Error: Orphaned: usr/lib64/libQt5Widgets.so.5
Error: Orphaned: usr/lib64/libQt5Widgets.so.5.15
Error: Orphaned: usr/lib64/libQt5Widgets.so.5.15.9
Error: Orphaned: usr/lib64/libQt5XcbQpa.so.5
Error: Orphaned: usr/lib64/libQt5XcbQpa.so.5.15
Error: Orphaned: usr/lib64/libQt5XcbQpa.so.5.15.9
Error: Orphaned: usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-emu-integration.so
Error: Orphaned: usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-kms-egldevice-integration.so
Error: Orphaned: usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-kms-integration.so
Error: Orphaned: usr/lib64/qt5/plugins/egldeviceintegrations/libqeglfs-x11-integration.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqevdevkeyboardplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqevdevmouseplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqevdevtabletplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqevdevtouchplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqlibinputplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/generic/libqtuiotouchplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/imageformats/libqgif.so
Error: Orphaned: usr/lib64/qt5/plugins/imageformats/libqico.so
Error: Orphaned: usr/lib64/qt5/plugins/imageformats/libqjpeg.so
Error: Orphaned: usr/lib64/qt5/plugins/platforminputcontexts/libcomposeplatforminputcontextplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/platforminputcontexts/libibusplatforminputcontextplugin.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqeglfs.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqlinuxfb.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqminimal.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqminimalegl.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqoffscreen.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqvnc.so
Error: Orphaned: usr/lib64/qt5/plugins/platforms/libqxcb.so
Error: Orphaned: usr/lib64/qt5/plugins/platformthemes/libqgtk3.so
Error: Orphaned: usr/lib64/qt5/plugins/platformthemes/libqxdgdesktopportal.so
Error: Orphaned: usr/lib64/qt5/plugins/printsupport/libcupsprintersupport.so
Error: Orphaned: usr/lib64/qt5/plugins/xcbglintegrations/libqxcb-egl-integration.so
Error: Orphaned: usr/lib64/qt5/plugins/xcbglintegrations/libqxcb-glx-integration.so

